# bash/files/profile.d/history.sh
export HISTSIZE=1000
export HISTIGNORE="&:[bf]g:exit"
export HISTCONTROL=ignoredups:ignorespace
shopt -s histappend

# check the window size after each command and, if necessary,
# update the values of LINES and COLUMNS.
shopt -s checkwinsize
